# launchd/com.sparkmonitor.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sparkmonitor.agent</string>

    <!-- Path to the installed binary. setup.sh installs to /Applications. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/SparkMonitor.app/Contents/MacOS/SparkMonitor</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Uncomment to override defaults:
    <key>EnvironmentVariables</key>
    <dict>
        <key>SPARK_HOST</key>
        <string>nvidia-dgx-spark</string>
        <key>SPARK_POLL_SECS</key>
        <string>15</string>
    </dict>
    -->

    <key>StandardErrorPath</key>
    <string>/tmp/sparkmonitor.err.log</string>
    <key>StandardOutPath</key>
    <string>/tmp/sparkmonitor.out.log</string>
</dict>
</plist>
